# docker/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: NSPanel Manager
# Runs the Nginx daemon
# ==============================================================================

# Wait for NSPanel Manager to become available
bashio::net.wait_for 8000

bashio::log.info "Starting NGINX..."

# nginx -T

exec nginx